# playbook/roles/openstack/templates/node-exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
ExecStart=/opt/tools/node_exporter
User=nobody

[Install]
WantedBy=multi-user.target
